# roles/application/templates/lighthouse-api.upstart.j2
start on startup
respawn

exec su - {{user}} -c 'cd /opt/baseline-api; export PORT=3000; export RAILS_ENV=production; rails server >> /var/log/baseline-api.log'
