# roles/monitoring/templates/riemann.upstart.j2
start on startup
respawn

exec su - {{user}} -c 'riemann'
